# agent/packaging/nodeguard-agent.service
[Unit]
Description=NodeGuard read-only Linux monitoring agent
Documentation=https://github.com/HackintoshMatrix7132/NodeGuard
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/nodeguard-agent run --config /etc/nodeguard-agent/config.json
StateDirectory=nodeguard-agent
StateDirectoryMode=0700
Restart=on-failure
RestartSec=10s
TimeoutStopSec=30s
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=-/var/lib/apt/lists -/var/cache/apt -/var/lib/apt/periodic
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
LockPersonality=true

[Install]
WantedBy=multi-user.target
